# deploy/systemd/sh-manager-operations@.service
# SPDX-FileCopyrightText: 2026 Humdek, University of Bern
# SPDX-License-Identifier: MPL-2.0
#
# Supervised CMS-update loop for ONE SelfHelp instance (systemd template unit).
# `%i` is the instance id. Enable one per instance:
#
#   sudo systemctl enable --now sh-manager-operations@<instance-id>.service
#
# It runs `sh-manager instance process-operations <id> --watch`, which claims and
# executes every CMS-requested update for that instance, writing live status back
# so a request never stays on "requested". Restart=always keeps it supervised.
#
# The default transport execs into the instance's backend container and
# authenticates with the token generated at install — no env file needed.
# Remote/HTTP setups put SELFHELP_MANAGER_TOKEN and
# SHM_OPERATIONS_EXTRA_ARGS=--backend-url http://... into the optional 0600
# EnvironmentFile (see deploy/README.md) so secrets never live in the unit or
# process list.

[Unit]
Description=SelfHelp Manager update-operations loop for instance %i
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
# Manager state root (inventory + per-instance manifests/secrets).
Environment=SELFHELP_MANAGER_ROOT=/var/lib/sh-manager
# Optional per-instance overrides (token + extra args for HTTP transport).
EnvironmentFile=-/etc/sh-manager/%i.env
ExecStart=/usr/bin/sh-manager --root ${SELFHELP_MANAGER_ROOT} instance process-operations %i \
  --watch --interval 15 $SHM_OPERATIONS_EXTRA_ARGS
Restart=always
RestartSec=10
# Hardening: the loop only needs the manager root + Docker socket.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
