# deploy/cron/sh-manager-operations.crontab.example
# SPDX-FileCopyrightText: 2026 Humdek, University of Bern
# SPDX-License-Identifier: MPL-2.0
#
# Cron alternative to the systemd template unit for hosts without systemd.
# Each line drains ALL pending CMS-requested operations for one instance every
# minute (process-operations without --watch is one-shot: claim+execute every
# pending op, then exit). Add one line per instance.
#
# Install for the manager service account:
#   crontab -u sh-manager deploy/cron/sh-manager-operations.crontab.example
#
# The default transport execs into the backend container with the token
# generated at install — the env file is optional. Remote/HTTP setups put
# SELFHELP_MANAGER_TOKEN and SHM_OPERATIONS_EXTRA_ARGS=--backend-url http://...
# into the 0600 env file (see deploy/README.md) so the token stays out of the
# crontab and process list. SHELL=/bin/bash is required for the `set -a; . file`.

SHELL=/bin/bash
SELFHELP_MANAGER_ROOT=/var/lib/sh-manager

# Replace <instance-id>. One line per instance.
* * * * * set -a; [ -f /etc/sh-manager/<instance-id>.env ] && . /etc/sh-manager/<instance-id>.env; set +a; /usr/bin/sh-manager --root "$SELFHELP_MANAGER_ROOT" instance process-operations <instance-id> $SHM_OPERATIONS_EXTRA_ARGS >> /var/log/sh-manager/<instance-id>-operations.log 2>&1
